# radxa-system-config-common/usr/lib/systemd/system/sddm.service.d/10-workaround-slow-shutdown.conf
# NOTE: https://github.com/sddm/sddm/issues/1476#issuecomment-1140615746

[Service]
TimeoutStopSec=5s
