# debian/radxa-system-config-disable-sleep.postinst
#!/bin/sh
set -e

case $1 in
    configure)
        deb-systemd-helper mask sleep.target suspend.target hibernate.target hybrid-sleep.target
        ;;
esac
    
#DEBHELPER#

exit 0
